# etc/supervisord/conf.d/gunicorn.conf
[program:gunicorn]
command=/usr/local/bin/gunicorn -b unix:/tmp/gunicorn.sock mass.wsgi
directory=/app
user=nobody
autostart=true
autorestart=true
redirect_stderr=true
